# scripts/tws.service
[Unit]
Description=TWS Quantitative Trading System
After=network.target

[Service]
# User running the script
User=ubuntu
WorkingDirectory=/home/ubuntu/tws
# Environment Variables (Or use 1Password 'op run')
# EnvironmentFile=/home/ubuntu/tws/.env

# Command
# Using 'op run' for 1Password integration
# ExecStart=/usr/local/bin/op run --env-file="./config.env.tpl" -- /usr/bin/python3 -m quant_system.main config.json

# Simple Mode (Env Vars in .env or system)
ExecStart=/usr/bin/python3 -m quant_system.main config.json

# Auto Restart
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
